# wcc-client@.socket
[Socket]
ListenDatagram=0.0.0.0:%i

[Install]
WantedBy=network.target
