# wcc-client@.service
[Unit]
Exec=/storage/bin/wcc-client
StandardInput=socket
StandardOutput=journal

[Install]
WantedBy=network.target
